# startup.service
[Unit]
Description=My custom startup script

[Service]
ExecStart=/etc/init.d/startup.sh

[Install]
WantedBy=multi-user.target